# config/celery.upstart
description "I Want Better Food Celeryd"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [016]
respawn


exec /var/www/signmeup/env/bin/python /var/www/signmeup/server/manage.py celeryd

